# lib/capops/templates/upstart/delayed_job.conf.erb
description "delayed job workers daemon"

start on runlevel [2345]
stop on runlevel [06]

respawn
respawn limit 10 90

script
  PATH=$PATH:/opt/earlydoc/current/vendor/cache
  RAILS_ENV=<% enviornment %>

  cd <%= current_path %>
  RAILS_ENV=<%= environment %> sudo bundle exec /opt/earlydoc/current/script/delayed_job run >> /opt/earlydoc/current/log/dj.log 2>&1
end script